# overlay/etc/cron.d/nextcloud
# run cron.php for nextcloud
5-59/15 * * * * www-data 	/usr/bin/php -f /usr/share/nextcloud/cron.php /dev/null 2>&1
