# overlay/etc/cron.d/nextcloud
# run cron.php for nextcloud
*/15 * * * * www-data 	/usr/bin/php -f /var/www/nextcloud/occ /dev/null 2>&1
